# raspberry-pi/emp-scanner.service
[Unit]
Description=EMP Access Scanner
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=root
WorkingDirectory=/opt/emp-scanner/raspberry-pi
ExecStart=/opt/emp-scanner/venv/bin/python -m emp_scanner.main
Restart=always
RestartSec=5
WatchdogSec=120
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
